# deployment/systemd/citations-backend.service
[Unit]
Description=Citation Format Checker Backend
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/citations
Environment=PATH=/opt/citations/venv/bin
ExecStart=/opt/citations/venv/bin/uvicorn backend.app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target